# extras/iniit-scripts/upstart/beanstalkd.conf
description "simple, fast work queue"

start on filesystem
stop on runlevel [!2345]

respawn
respawn limit 5 2

exec su nobody -c 'exec /usr/bin/beanstalkd'
