# extras/iniit-scripts/systemd/beanstalkd.service
[Unit]
Description=Beanstalkd is a simple, fast work queue

[Service]
User=nobody
ExecStart=/usr/bin/beanstalkd
