# extras/iniit-scripts/systemd/beanstalkd.socket
[Unit]
Description=Socket for beanstalkd, a simple, fast work queue

[Socket]
ListenStream=11300

[Install]
WantedBy=sockets.target
