# adm/upstart/beanstalkd.conf
description "simple, fast work queue"

start on filesystem
stop on runlevel [!2345]

respawn
respawn limit 5 2

setuid nobody

exec /usr/bin/beanstalkd
